# templates/voidbr-samba/pkg/etc/profile.d/voidbr-samba.sh
export PATH=/opt/samba/bin:/opt/samba/sbin:$PATH
